# frameworks/Python/web2py/web2py/applications/admin/cron/crontab
10	*	*	*	*	root **applications/admin/cron/expire_sessions.py

